# crontab
@reboot sh /home/pi/studio/launcher.sh >/home/pi/logs/cronlog 2>&1
